# conf/supervisor/glances.conf
[program:glances]
command=glances -w
directory=/home/ubuntu
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true